# dist/watchdog.service
[Unit]
Description=Watchdog Daemon

[Service]
Type=simple
ExecStart=/usr/bin/watchdog

[Install]
WantedBy=multi-user.target
